# roles/unattended-upgrades/files/systemd/system/apt-daily-upgrade.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=05:20
RandomizedDelaySec=0
