# roles/auto-upgrade/files/systemd/system/apt-daily.timer.d/override.conf
[Timer]
OnCalendar=
OnCalendar=05:00
RandomizedDelaySec=0
